# templates/default/gdash_init.erb
#!/bin/bash
#
# Generated by Chef for <%= node[:fqdn] %> (gdash::default)
# Do not edit locally
#
# gdash: a dashboard for graphite
#
# chkconfig:   2345 80 20
# description: A simple Graphite dashboard built using Twitter's Bootstrap. \
#              Project page: https://github.com/ripienaar/gdash

# Source CentOS function library
. /etc/rc.d/init.d/functions

# Name of the service
NAME=gdash

# Runas user & path to gdash application
USER=<%= node.gdash.owner %>
APP_PATH=<%= node.gdash.base %>

# Gdash pid & lockfiles
PID=$APP_PATH/$NAME.pid
LOCK=/var/lock/subsys/$NAME

# Gdash runtime options to pass at startup
OPTS="-c /etc/unicorn/$NAME.app -E <%= node.chef_environment %> -D"

# Ruby path update. Shouldn't need thanks to Bundler binstubs
RUBY_PATH="PATH=$PATH"

start() {
  cd $APP_PATH

  echo -n $"Starting $NAME: "
  daemon --check=$NAME --pidfile=$PID --user=$USER "$RUBY_PATH bundle exec unicorn $OPTS"
  retval=$?
  [ $retval -eq 0 ] && touch $LOCK
  echo

  return $retval
}

stop() {
  cd $APP_PATH

  echo -n $"Stopping $NAME: "
  killproc -p $PID
  retval=$?
  [ $retval -eq 0 ] && rm -f $LOCK
  echo

  return $retval
}

restart() {
  stop
  start
}

get_status() {
  status -p $PID $NAME
}

query_status() {
  get_status >/dev/null 2>&1
}

case "$1" in
  start)
    query_status && exit 0
    start
    ;;
  stop)
    query_status || exit 0
    stop
    ;;
  restart)
    restart
    ;;
  status)
    get_status
    ;;
  *)
    N=/etc/init.d/$NAME
    echo "Usage: $N {start|stop|restart|status}" >&2
    exit 1
    ;;
esac

exit 0
